# supervisor.conf
[program:studynotes]
user=feross
process_name=studynotes
directory=/home/feross/www/studynotes.org/
command=/usr/bin/node index.js
environment=NODE_ENV=production
startsecs=10
startretries=50
redirect_stderr=true
stdout_logfile=/home/feross/www/log/studynotes.org/access.log